# systemd/archility-weekly.service
[Unit]
Description=Archility weekly architecture audit and diagram render
Documentation=https://github.com/casonk/traction-control
After=network.target

[Service]
Type=oneshot

# Override these via a drop-in or environment file as needed.
Environment=PORTFOLIO_ROOT=/mnt/4tb-m2/git
Environment=ARCHILITY_CMD=/home/user/.local/bin/archility
Environment=LOG_DIR=%h/.local/share/archility-weekly

ExecStart=/mnt/4tb-m2/git/util-repos/traction-control/scripts/archility-weekly.sh

# Keep the last 3 journal entries for this unit visible with:
#   journalctl --user -u archility-weekly.service
StandardOutput=journal
StandardError=journal
